# docker/app/supervisor.conf
[supervisord]
nodaemon=true

[program:app-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/app/artisan queue:work --sleep=3 --tries=3
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
user=root
stdout_logfile=/var/www/app/storage/logs/worker.log

[program:segmentation_socket_server_23456]
command=python3 /var/www/app/app/Python/decimer_segmentation_server.py 23456
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/segmentation_socket_server_23456.log

[programm:segmentation_socket_server_23457]
command=python3 /var/www/app/app/Python/decimer_segmentation_server.py 23457
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/segmentation_socket_server_23457.log

[programm:segmentation_socket_server_23458]
command=python3 /var/www/app/app/Python/decimer_segmentation_server.py 23458
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/segmentation_socket_server_23458.log

[programm:segmentation_socket_server_23459]
command=python3 /var/www/app/app/Python/decimer_segmentation_server.py 23459
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/segmentation_socket_server_23459.log

[programm:segmentation_socket_server_23460]
command=python3 /var/www/app/app/Python/decimer_segmentation_server.py 23460
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/segmentation_socket_server_23460.log

[programm:segmentation_socket_server_23461]
command=python3 /var/www/app/app/Python/decimer_segmentation_server.py 23461
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/segmentation_socket_server_23461.log

[programm:segmentation_socket_server_23462]
command=python3 /var/www/app/app/Python/decimer_segmentation_server.py 23462
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/segmentation_socket_server_23462.log

[programm:segmentation_socket_server_23463]
command=python3 /var/www/app/app/Python/decimer_segmentation_server.py 23463
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/segmentation_socket_server_23463.log

[programm:segmentation_socket_server_23464]
command=python3 /var/www/app/app/Python/decimer_segmentation_server.py 23464
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/segmentation_socket_server_23464.log

[programm:segmentation_socket_server_23465]
command=python3 /var/www/app/app/Python/decimer_segmentation_server.py 23465
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/segmentation_socket_server_23465.log


[program:ocsr_socket_server_65432]
command=python3 /var/www/app/app/Python/decimer_predictor_server.py 65432
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/ocsr_socket_server_65432.log

[program:ocsr_socket_server_65433]
command=python3 /var/www/app/app/Python/decimer_predictor_server.py 65433
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/ocsr_socket_server_65433.log

[programm:ocsr_socket_server_65434]
command=python3 /var/www/app/app/Python/decimer_predictor_server.py 65434
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/ocsr_socket_server_65434.log

[programm:ocsr_socket_server_65435]
command=python3 /var/www/app/app/Python/decimer_predictor_server.py 65435
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/ocsr_socket_server_65435.log

[programm:ocsr_socket_server_65436]
command=python3 /var/www/app/app/Python/decimer_predictor_server.py 65436
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/ocsr_socket_server_65436.log

[programm:ocsr_socket_server_65437]
command=python3 /var/www/app/app/Python/decimer_predictor_server.py 65437
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/ocsr_socket_server_65437.log

[program:stout_socket_server_12345]
command=python3 /var/www/app/app/Python/stout_predictor_server.py 12345
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/stout_socket_server_12345.log

[programm:stout_socket_server_12346]
command=python3 /var/www/app/app/Python/stout_predictor_server.py 12346
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/stout_socket_server_12346.log

[programm:stout_socket_server_12347]
command=python3 /var/www/app/app/Python/stout_predictor_server.py 12347
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/stout_socket_server_12347.log

[programm:stout_socket_server_12348]
command=python3 /var/www/app/app/Python/stout_predictor_server.py 12348
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/stout_socket_server_12348.log

[programm:stout_socket_server_12349]
command=python3 /var/www/app/app/Python/stout_predictor_server.py 12349
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/stout_socket_server_12349.log

[programm:stout_socket_server_12350]
command=python3 /var/www/app/app/Python/stout_predictor_server.py 12350
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/stout_socket_server_12350.log

[program:clear_config]
process_name=clear_config
command=php /var/www/app/artisan config:clear
autorestart=false
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/clear_config.log

[program:clear_cache]
process_name=clear_cache
command=php /var/www/app/artisan cache:clear
autorestart=false
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/clear_cache.log

[program:laravel_link]
process_name=create_storage_link
command=php /var/www/app/artisan storage:link
autorestart=false
redirect_stderr=true
stdout_logfile=/var/www/app/storage/logs/storage_link.log